# apirouter.service
[Unit]
Description=API Router - Multi-upstream proxy with failover
Documentation=https://github.com/winterhuan/api-router
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/maiscrm/workspace/apirouter
ExecStart=/home/maiscrm/workspace/apirouter/target/release/apirouter --port 1999
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
LimitNOFILE=65536
Environment=RUST_LOG=info
Environment=HTTP_PROXY=http://127.0.0.1:10808
Environment=HTTPS_PROXY=http://127.0.0.1:10808
Environment=NO_PROXY=localhost,127.0.0.1,::1,evomap.ai,.evomap.ai,.feishu.cn,.larksuite.com,.larkoffice.com,.feishucdn.com,.bytedance.com,.volces.com,.volcengine.com,.aliyuncs.com,.qwen.ai,.minimax.io,.minimaxi.com,.moonshot.ai,.baidu.com,.baidubce.com,.deepseek.com,.nvidia.com,longcat.chat,.longcat.chat

[Install]
WantedBy=multi-user.target
