# apirouter.service
[Unit]
Description=API Router - Multi-upstream proxy with failover
Documentation=https://github.com/winterhuan/api-router
After=network.target

[Service]
Type=simple
# User 和 WorkingDirectory 由安装脚本动态设置
WorkingDirectory=%h/workspace/apirouter
ExecStart=%h/workspace/apirouter/target/release/apirouter --port 1999
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
LimitNOFILE=65536
Environment=RUST_LOG=info

[Install]
WantedBy=default.target
